# python-ai/supervisor-ai-service.conf
# Supervisor配置文件 - Python AI服务
# 功能：进程守护、自动重启、日志管理
# 需求：5.4, 5.5
#
# 安装说明：
# 1. 安装supervisor: sudo apt-get install supervisor (Ubuntu/Debian)
#                    或 sudo yum install supervisor (CentOS/RHEL)
# 2. 复制此文件到: /etc/supervisor/conf.d/ai-service.conf
# 3. 修改路径配置（directory, command, stdout_logfile, stderr_logfile）
# 4. 重新加载配置: sudo supervisorctl reread
# 5. 启动服务: sudo supervisorctl update
# 6. 查看状态: sudo supervisorctl status ai-service
#
# 常用命令：
# - 启动服务: sudo supervisorctl start ai-service
# - 停止服务: sudo supervisorctl stop ai-service
# - 重启服务: sudo supervisorctl restart ai-service
# - 查看日志: sudo supervisorctl tail -f ai-service

[program:ai-service]
# 服务名称
command=python3 app.py
# 工作目录（需要修改为实际路径）
directory=/path/to/python-ai
# 运行用户（建议使用非root用户）
user=www-data
# 进程数量
numprocs=1
# 自动启动
autostart=true
# 自动重启
autorestart=true
# 启动等待时间（秒）
startsecs=10
# 停止等待时间（秒）
stopwaitsecs=10
# 停止信号
stopsignal=TERM
# 重启尝试次数
startretries=3
# 重定向stderr到stdout
redirect_stderr=true
# 标准输出日志文件（需要修改为实际路径）
stdout_logfile=/path/to/python-ai/logs/supervisor-ai-service.log
# 日志文件最大大小
stdout_logfile_maxbytes=50MB
# 日志文件备份数量
stdout_logfile_backups=10
# 标准错误日志文件（需要修改为实际路径）
stderr_logfile=/path/to/python-ai/logs/supervisor-ai-service-error.log
# 错误日志文件最大大小
stderr_logfile_maxbytes=50MB
# 错误日志文件备份数量
stderr_logfile_backups=10
# 环境变量
environment=PORT="5000",GRPC_PORT="50051",PYTHONUNBUFFERED="1"
# 进程优先级
priority=999
